# deploy/systemd/severino-hq-secrets.service
[Unit]
Description=Refresh Severino HQ secrets from 1Password
Documentation=https://developer.1password.com/docs/cli/secrets-scripts/
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
LoadCredentialEncrypted=op_service_account_token
Environment=OP_CONFIG_DIR=/run/severino-hq-op
RuntimeDirectory=severino-hq-op
RuntimeDirectoryMode=0700
ExecStart=/opt/apps/severino-hq/scripts/refresh-secrets.sh

NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/opt/apps/severino-hq/secrets
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
