# deploy/systemd/severino-hq-secrets.timer
[Unit]
Description=Periodically refresh Severino HQ secrets from 1Password

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=yes
Unit=severino-hq-secrets.service

[Install]
WantedBy=timers.target

